# systemd/system/network-connected.target
[Unit]
Description=Network is connected and able to send packets
After=network-online.target
