# systemd/system/reflector.timer
[Unit]
Description=Runs Reflector after one minute

[Timer]
OnBootSec=1min
Unit=reflector.service

[Install]
WantedBy=multi-user.target
